# ansible/roles/desktop/templates/atomik-appnap.service.j2
[Unit]
Description=Atomik App Nap
After=graphical-session.target

[Service]
ExecStart=pkexec /usr/local/lib/atomik/atomik-appnap.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
